# ipkg/mjpg-streamer/postrm
#!/bin/sh
/usr/sbin/update-rc.d -s mjpg-streamer remove
